# httpd/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
;user=apache
;directory=/var/www/html
